# files/systemd/system/gm-fdo-accumulate.service
[Unit]
Description=Golden Master profile accumulator (llvm-profgen)

[Service]
Type=oneshot
ExecStart=/usr/local/bin/gm-fdo-accumulate
User=root
