# files/systemd/system/gm-perf-profiler.service
[Unit]
Description=Golden Master process-aware perf sampler
After=multi-user.target

[Service]
Type=simple
ExecStart=/usr/local/bin/gm-perf-profiler
Restart=always
RestartSec=3

# perf needs elevated capabilities; simplest is root
User=root

[Install]
WantedBy=multi-user.target
